# infra/systemd/dashmoney-update-prices.service
[Unit]
Description=DashMoney - Update daily prices (one-shot)
After=network.target dashmoney-backend.service
Wants=dashmoney-backend.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dashmoney-update-prices.sh
