# infra/systemd/dashmoney-update-prices.timer
[Unit]
Description=DashMoney - Schedule daily price update

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
